# environment/RunIPythonOnStartup.conf
description "Run IPython on startup"
start on startup
task
cd ~
sudo ipython notebook --profile=nbserver